# rpi/alpine-mark-i/overlays/mark-i/etc/local.d/0.bootstrap.start
#!/bin/sh

apk update
apk add

update-ca-certificates
ntpd -s

chown -R root:root /etc
